# modules/openvpn/files/etc/cron.d/cert_check
0 12 * * *	root	/usr/bin/cert_check -d 10
